# snat-reflector.service
[Unit]
Description=SNAT Reflactor

[Service]
Type=exec
ExecStart=/usr/local/bin/snat-reflector

[Install]
WantedBy=multi-user.target
